# package/lean/luci-app-sysmonitor/root/etc/hotplug.d/iface/90-sysmonitor
[ $INTERFACE = lan6 ] && touch /tmp/ipv6.sign
[ $INTERFACE = lan ] && touch /tmp/ipv4.sign

